# app-service/supervisord.conf
[supervisord]
nodaemon=true

[program:django]
command=python3 /opt/app/manage.py runserver 0.0.0.0:80